# service/renderfarm-server.service
[Unit]
Description=RenderFarm Server
After=network.target

[Service]
Type=simple
WorkingDirectory={{INSTALL_DIR}}/server
ExecStart=/usr/bin/uv run python main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
